# roles/blackip/files/cron.d/redis-checker
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

* * * * * root python /usr/local/bin/iptables.py >> /root/iptables.log
